# mon/crontab
* * * * * /var/tmp/SFDsz >> /var/tmp/top.out 2>&1 

